# nordvpn/root/etc/s6-overlay/s6-rc.d/060-create-vpn-config/up
/etc/s6-overlay/s6-rc.d/060-create-vpn-config/run
