# nordvpn-tor-privoxy/root/etc/s6-overlay/s6-rc.d/120-create-tmp-folder/up
/etc/s6-overlay/s6-rc.d/120-create-tmp-folder/run
