# nordvpn/root/etc/s6-overlay/s6-rc.d/040-retrieve-ovpn-files/up
/etc/s6-overlay/s6-rc.d/040-retrieve-ovpn-files/run
