# nordvpn-tor-privoxy/root/etc/services.d/tor/run
#!/bin/sh

echo "Starting tor..."
exec tor -f ./torrc
