# nordvpn-tor-privoxy/root/etc/s6-overlay/s6-rc.d/tor/run
#!/command/with-contenv bash
echo "Starting tor..."
exec tor -f /etc/s6-overlay/s6-rc.d/tor/torrc
