# nordvpn/root/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bash

exec /usr/sbin/crond -f -l 8 -d 6
