# python-alpine-s6/root/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/crond -f -l 8 -d 6
