# nordvpn/root/etc/s6-overlay/s6-rc.d/020-iptables/up
/etc/s6-overlay/s6-rc.d/020-iptables/run
