# nordvpn-privoxy/root/etc/services.d/privoxy/run
#!/bin/sh
echo "Starting privoxy..."
exec /usr/sbin/privoxy --no-daemon /etc/privoxy/config
